# supervisor.conf
[program:application]
directory=/srv
command=python3 /srv/application.py
environment=DB_CONNECTION=postgresql://demo:password@localhost:5432/demo